# usr/share/odroid-kiosk/odroid-kiosk-weston.service.in
[Unit]
Description=Qt Kiosk Browser
After=weston.service
Requires=weston.service

[Service]
User=@@DEFAULT_USER@@
EnvironmentFile=-/etc/default/odroid-kiosk
Restart=yes
Type=simple
ExecStart=/usr/share/odroid-kiosk/weston-autostart.sh

[Install]
WantedBy=multi-user.target weston.service
