# usr/share/odroid-kiosk/weston.service.in
[Unit]
Description=Weston Wayland Compositor (on tty7)
RequiresMountsFor=/run
Conflicts=getty@tty7.service
After=systemd-user-sessions.service getty@tty7.service

[Service]
User=@@DEFAULT_USER@@
PermissionsStartOnly=true
EnvironmentFile=-/etc/default/odroid-kiosk

# Log us in via PAM so we get our XDG & co. environment and
# are treated as logged in so we can use the tty:
PAMName=login

# Grab tty7
UtmpIdentifier=tty7
TTYPath=/dev/tty7
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# stderr to journal so our logging doesn't get thrown into /dev/null
StandardOutput=tty
StandardInput=tty
StandardError=journal

# Weston does not successfully change VT, nor does systemd place us on
# the VT it just activated for us. Switch manually:
ExecStartPre=/bin/chvt 7
ExecStartPre=/usr/bin/sudo /bin/mkdir -p /run/user/@@DEFAULT_USERID@@
ExecStartPre=/usr/bin/sudo /bin/chown @@DEFAULT_USER@@:@@DEFAULT_USER@@ /run/user/@@DEFAULT_USERID@@
ExecStartPre=/bin/chmod 0700 /run/user/@@DEFAULT_USERID@@
ExecStart=/usr/bin/weston-launch -- --xwayland $OPTARGS
IgnoreSIGPIPE=no

[Install]
WantedBy=multi-user.target
