# debian/odroid-kiosk-weston.postrm
#!/bin/sh

for service in odroid-kiosk-weston weston; do
	systemctl stop $service
	systemctl disable $service
done

chvt 1

#DEBHELPER#
