# debian/odroid-kiosk-weston.postrm
#!/bin/sh

service=odroid-kiosk-weston

if [ -f /etc/default/odroid-kiosk-weston ]; then
	. /etc/default/odroid-kiosk-weston
fi

systemctl stop $service
systemctl disable service
rm -f /usr/lib/systemd/system/$service

chvt 1

#DEBHELPER#
